# src/Coalesce.Dist/coalesce-karaf-dist/supervisord.conf
[supervisord]
nodaemon=true
  
[program:karaf]
command=/opt/karaf/bin/karaf
