# src/Coalesce.Dist/coalesce-karaf-dist/supervisord.conf
[supervisord]
nodaemon=true
  
[program:karaf]
command=/bin/sh -c "rm -f $KARAF_HOME/instances/instance.properties && karaf start"
