# src/Coalesce.Dist/coalesce-karaf-dist/supervisord.conf
[supervisord]
nodaemon=true
  
[program:karaf]
command=/bin/sh -c "rm -f /opt/karaf/instances/instance.properties && /opt/karaf/bin/karaf clean"
